# biscavolley-bot.service
[Unit]
Description=BiscaVolley Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/biscavolley
ExecStart=/opt/biscavolley/venv/bin/python3 /opt/biscavolley/biscavolley-bot.py
Restart=always
RestartSec=5
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/biscavolley/environment

[Install]
WantedBy=multi-user.target

